# scstack/templates/service/scstack-service.conf.erb
# SCStack service

description     "SCStack Service"
author          "Patxi Gortázar <patxi.gortazar@gmail.com>"

# Change this path to your own value
env SCSTACK_DIR=<%=installFolder%>/scstack-service

# Don't change the following lines
##############################################
start on (net-device-up
          and local-filesystems
	      and runlevel [2345])
stop on runlevel [016]

respawn
respawn limit 15 5

script
chdir $SCSTACK_DIR
exec /etc/init.d/scstack-service > $SCSTACK_DIR/scstack-service.log 2>&1
end script
